# service/sg2042-pisces/powergood/powergood.service
[Unit]
Description=Power Good Service
After=multi-user.target

[Service]
ExecStart=/usr/local/power_good.sh 
[Install]
WantedBy=default.target
 
